# gunicorn.service
[Unit]
Description=Gunicorn daemon for Ingush Translate Django application
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/ingushtranslate
Environment="PATH=/var/www/ingushtranslate/venv/bin"
Environment="DJANGO_SETTINGS_MODULE=ingushtranslate.settings"
Environment="DEBUG=False"
Environment="ALLOWED_HOSTS=localhost,127.0.0.1,your-domain.com"
ExecStart=/var/www/ingushtranslate/venv/bin/gunicorn --workers 3 --bind unix:/var/www/ingushtranslate/ingushtranslate.sock ingushtranslate.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true

[Install]
WantedBy=multi-user.target 